# bare-silabs-multiprotocol/rootfs/etc/services.d/zigbeed/finish
#!/bin/bash
#==============================================================================
# EmberZNet Zigbee Daemon finish script
#==============================================================================
echo "zigbeed ended with exit code ${1} (signal ${2})..."
